# deployment/systemd/pivision-worker.service
[Unit]
Description=PiVision background worker
After=network.target pivision-server.service
Wants=pivision-server.service

[Service]
Type=simple
User=djg
WorkingDirectory=/home/djg/Projects/PiVision
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/djg/Projects/PiVision/.venv/bin/python -m backend.worker
Restart=always
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pivision-worker
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
