# deployment/systemd/pivision-server.service
[Unit]
Description=PiVision API server
After=network.target

[Service]
Type=simple
User=djg
WorkingDirectory=/home/djg/Projects/PiVision
Environment=PORT=8080
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/djg/Projects/PiVision/.venv/bin/python -m backend.server
Restart=always
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pivision-server
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
